# cron.d/backupkg.cron
0 22 * * * /home/matt/resources/mygit/myarch/cron.d/backupkg.sh
