# cron.d/matt.cron
0 12,22 * * * /home/matt/resources/mygit/myarch/bin/updatearchlinuxiso.sh
